# simple.conf
[program:myserver]
command=gunicorn website:app -w 2 -b 0.0.0.0:80

[supervisord]
logfile=/home/ubuntu/supervisord.log
loglevel=debug
user=root